# borne/bin/systemd/beequest.service
[Unit]
Description=Démarrage de la borne
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=turbopollen
WorkingDirectory=/home/turbopollen/en-quete-abeilles/borne
ExecStart=/home/turbopollen/en-quete-abeilles/borne/bin/start.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target